# files/.config/autostart/synergy.desktop
[Desktop Entry]
Name=Synergy
Type=Application
Exec=synergy
Terminal=false