# etc/profile.d/custom-aliases.sh
#!/usr/bin/env bash

mkdircd () {
    mkdir "$1" && cd "$1"
}
